# supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:sensu-client]
command=/opt/sensu/bin/sensu-client
